# test/fixtures/extractor-repo/ops/panels.service
[Unit]
Description=Panels worker
After=network.target

[Service]
ExecStart=/usr/bin/panels-worker

[Install]
WantedBy=multi-user.target
